# Nezha/deploy/nezha.service.example
# Copier vers /etc/systemd/system/nezha.service puis :
#   systemctl daemon-reload && systemctl enable --now nezha
#
# Plus fiable que « npm run start » sous systemd (charge .env comme en shell).

[Unit]
Description=Nezha Medical (Next.js)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
WorkingDirectory=/var/www/nezha/Nezha
Environment=NODE_ENV=production
Environment=PORT=3001
ExecStart=/bin/bash -lc 'set -a && source /var/www/nezha/Nezha/.env && set +a && exec /usr/bin/node /var/www/nezha/Nezha/node_modules/next/dist/bin/next start -p 3001'
Restart=always
RestartSec=8
TimeoutStartSec=120
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
